# deploy/aws-ec2/brr-api.service.example
# /etc/systemd/system/brr-api.service
#
# Runs the BRR Liquor Soft api-server on boot and restarts it on failure.
# Reads its environment from /etc/brr/brr-api.env (mode 600, owned by root).
#
# Install:
#   sudo cp deploy/aws-ec2/brr-api.service.example /etc/systemd/system/brr-api.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now brr-api.service
#
# Inspect:
#   systemctl status brr-api
#   journalctl -u brr-api -f

[Unit]
Description=BRR Liquor Soft API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Where the deployed code lives. Adjust if you put it elsewhere.
# Expected layout (created by the runbook):
#   /opt/brr/repo/                                 <- git checkout
#   /opt/brr/repo/artifacts/api-server/node_modules <- runtime deps (persistent
#                                                     across rebuilds; populated
#                                                     by `pnpm install --prod
#                                                     --filter @workspace/api-server`)
#   /opt/brr/repo/release/api/dist/index.mjs       <- esbuild bundle (rewritten
#                                                     by every release build)
#
# WorkingDirectory points at the source tree so node's module resolver finds
# the persistent node_modules. ExecStart points at the bundle, which is what
# actually changes between releases.
WorkingDirectory=/opt/brr/repo/artifacts/api-server
ExecStart=/usr/bin/node --enable-source-maps /opt/brr/repo/release/api/dist/index.mjs

EnvironmentFile=/etc/brr/brr-api.env

# Run as an unprivileged user that owns /opt/brr.
User=brr
Group=brr

Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Light hardening. Tighten further as you like.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
